# workflows/agenda/agenda-service.service
[Unit]
Description=BLBoardSolutions Agenda docx-Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/scripts
ExecStart=/usr/bin/node /opt/scripts/agenda-service.js
Restart=on-failure
RestartSec=3
# Nur lokal erreichbar; optional Token setzen:
# 0.0.0.0, damit der n8n-Docker-Container den Dienst über 172.17.0.1 erreicht.
# Port 3001 unbedingt per Firewall nach außen sperren: ufw deny 3001
Environment=AGENDA_HOST=0.0.0.0
Environment=AGENDA_PORT=3001
# Environment=AGENDA_TOKEN=dein-geheimes-token
Environment=AGENDA_TEMPLATE=/opt/scripts/agenda_template.docx
Environment=AGENDA_OUTPUT_DIR=/opt/scripts/output
User=root

[Install]
WantedBy=multi-user.target
